# config/root/etc/services.d/openresty/run
#!/usr/bin/with-contenv bash

exec /usr/bin/openresty -g 'daemon off;'
